# config/cron/sla_monitor.cron
*/5 * * * * cd /var/www/msp_manager && bundle exec rails runner "Sla::Monitor.run"
